# deploy/debian/emmisuche-reindex.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=root

0 1 * * * emmisuche /usr/local/sbin/emmisuche-reindex
